# scripts/com.bhagya.mcp-kubernetes.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bhagya.mcp-kubernetes</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/bhagya/go/bin/mcp-kubernetes</string>
        <!-- Add additional arguments here if needed -->
        <!-- <string>-port</string> -->
        <!-- <string>8080</string> -->
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/bhagya/Library/Logs/mcp-kubernetes.err</string>
    <key>StandardOutPath</key>
    <string>/Users/bhagya/Library/Logs/mcp-kubernetes.log</string>
</dict>
</plist>
